# deploy/runtime/supervisor/ccdash.conf
[group:ccdash]
programs=ccdash-api,ccdash-worker,ccdash-frontend
priority=999

[program:ccdash-api]
directory=/opt/ccdash
command=/bin/sh -lc '. /etc/ccdash/ccdash-shared.env && . /etc/ccdash/ccdash-api.env && export PYTHONPATH=/opt/ccdash && exec /opt/ccdash/backend/.venv/bin/python -m uvicorn backend.runtime.bootstrap_api:app --host 0.0.0.0 --port 8000'
user=ccdash
autostart=true
autorestart=true
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/ccdash/api.log
stderr_logfile=/var/log/ccdash/api.err.log

[program:ccdash-worker]
directory=/opt/ccdash
command=/bin/sh -lc '. /etc/ccdash/ccdash-shared.env && . /etc/ccdash/ccdash-worker.env && export PYTHONPATH=/opt/ccdash && exec /opt/ccdash/backend/.venv/bin/python -m backend.worker'
user=ccdash
autostart=true
autorestart=true
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/ccdash/worker.log
stderr_logfile=/var/log/ccdash/worker.err.log

[program:ccdash-frontend]
directory=/opt/ccdash
command=/usr/bin/env npm run start:frontend
user=ccdash
environment=NODE_ENV="production"
autostart=true
autorestart=true
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/ccdash/frontend.log
stderr_logfile=/var/log/ccdash/frontend.err.log
