# deploy/runtime/systemd/ccdash-frontend.service
[Unit]
Description=CCDash frontend preview runtime
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ccdash
Group=ccdash
WorkingDirectory=/opt/ccdash
Environment=NODE_ENV=production
ExecStart=/usr/bin/env npm run start:frontend
Restart=on-failure
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
