# deploy/runtime/systemd/ccdash-api.service
[Unit]
Description=CCDash API runtime
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ccdash
Group=ccdash
WorkingDirectory=/opt/ccdash
EnvironmentFile=/etc/ccdash/ccdash-shared.env
EnvironmentFile=/etc/ccdash/ccdash-api.env
ExecStart=/opt/ccdash/backend/.venv/bin/python -m uvicorn backend.runtime.bootstrap_api:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
